# deploy/host/mcp-reconciler.service
# The deploy reconciler: applies chat-approved tool deploys (deploy/host/reconcile.py).
# Placeholders are stamped at install time -- see deploy/host/README.md:
#   sudo python3 deploy/host/reconcile.py --init --repo $(pwd) --user $(whoami)
#   sed "s|__REPO__|$(pwd)|; s|__RUN_AS__|$(whoami)|" deploy/host/mcp-reconciler.service \
#     | sudo tee /etc/systemd/system/mcp-reconciler.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now mcp-reconciler
[Unit]
Description=mcp-tools deploy reconciler (chat-approved tool deploys)
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/bin/python3 __REPO__/deploy/host/reconcile.py --watch --repo __REPO__
# Must be able to run `docker compose` (docker group or root) and write the repo.
User=__RUN_AS__
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
