# tools/hamilton/systemd/mcp-hamilton.service
[Unit]
# Hardened systemd SYSTEM unit template for an mcp-tools server.
# Placeholders: mcp-hamilton (unit + state dir name), 8064 (loopback MCP port),
# /home/wes/Documents/mcp-tools/tools/hamilton (tool dir holding server.py and .venv), 8076 (this tool's squid
# egress listener), wes (unprivileged run-as user).
# SYSTEM unit (not --user) so (a) IPAddressDeny/Allow actually enforce the egress
# wall and (b) it starts at boot. Install via scripts/install-system.sh (sudo once).
Description=MCP tool: mcp-hamilton (public, Google-OAuth gated via Cloudflare Tunnel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wes
Group=wes
WorkingDirectory=/home/wes/Documents/mcp-tools/tools/hamilton
# MCP_PORT is read by main() BEFORE .env loads, so it MUST live here, not in .env.
Environment=MCP_HOST=127.0.0.1
Environment=MCP_PORT=8064
Environment=HOME=/home/wes
Environment=PATH=/home/wes/.local/bin:/usr/local/bin:/usr/bin:/bin
# OAuth state stays writable under ProtectHome=read-only via StateDirectory (/var/lib/mcp-hamilton).
StateDirectory=mcp-hamilton
Environment=FASTMCP_HOME=%S/mcp-hamilton
# numba JIT cache must land somewhere writable: the venv lives in $HOME (read-only under
# ProtectHome), so pandas_ta's cache=True functions otherwise fail at import with "no
# locator available for file '.../pandas_ta/utils/_math.py'". Point numba at the
# StateDirectory. Must be in the UNIT (not .env): pandas_ta is imported before .env loads.
Environment=NUMBA_CACHE_DIR=%S/mcp-hamilton/numba-cache

# --- L2 egress: all outbound traffic MUST go through the egress proxy ---
# This tool's per-tool squid listener (:8076) enforces its domain
# allowlist (security/egress-proxy/allowlist/mcp-hamilton.txt). The IPAddress wall makes
# the proxy the ONLY route off-box; NO_PROXY keeps loopback (guardrail, resolver) direct.
Environment=HTTPS_PROXY=http://127.0.0.1:8076
Environment=HTTP_PROXY=http://127.0.0.1:8076
Environment=NO_PROXY=127.0.0.1,localhost
IPAddressDeny=any
IPAddressAllow=127.0.0.1/8

ExecStart=/home/wes/Documents/mcp-tools/tools/hamilton/.venv/bin/python server.py
Restart=on-failure
RestartSec=5

# --- Hardening (L3) ---
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# Enforceable as a SYSTEM unit (these silently failed in --user mode):
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectProc=invisible
ProcSubset=pid

[Install]
WantedBy=multi-user.target
