# tools/xmcp/systemd/mcp-xmcp.service
[Unit]
# Public, Google-OAuth-gated x-mcp instance for the Claude apps (desktop/web/mobile),
# served over the Cloudflare Tunnel at https://xmcp.secure-agentic-engineering.com.
# SYSTEM unit (not --user): runs as User=wes but is managed by the system manager so
# (a) IPAddressDeny/Allow actually enforce (the egress wall), (b) it starts at BOOT
# without login/linger. Install via scripts/install-system.sh (needs sudo once).
Description=MCP tool: x-mcp (public, Google-OAuth gated via Cloudflare Tunnel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wes
Group=wes
WorkingDirectory=/home/wes/Documents/mcp-tools/tools/xmcp
# MCP_PORT is read by main() BEFORE .env loads, so it MUST live here, not in .env.
Environment=MCP_HOST=127.0.0.1
Environment=MCP_PORT=8061
Environment=HOME=/home/wes
Environment=PATH=/home/wes/.local/bin:/usr/local/bin:/usr/bin:/bin
# State (OAuth client registrations + encrypted tokens) stays writable under
# ProtectHome=read-only via StateDirectory -> /var/lib/mcp-xmcp.
StateDirectory=mcp-xmcp
Environment=FASTMCP_HOME=%S/mcp-xmcp

# --- L2 egress: all outbound traffic MUST go through the egress proxy ---
# x-mcp's per-tool squid listener (:8073) enforces its domain allowlist. The
# IPAddress wall below makes the proxy the ONLY route off-box (a malicious dep that
# ignores HTTPS_PROXY hits the kernel drop). NO_PROXY keeps loopback (guardrail
# :8071, resolver) direct.
Environment=HTTPS_PROXY=http://127.0.0.1:8073
Environment=HTTP_PROXY=http://127.0.0.1:8073
Environment=NO_PROXY=127.0.0.1,localhost
# Only loopback is reachable directly: the proxy (:8073), guardrail (:8071), and the
# systemd-resolved stub (127.0.0.53). Everything else is dropped by the kernel.
IPAddressDeny=any
IPAddressAllow=127.0.0.1/8

ExecStart=/home/wes/Documents/mcp-tools/tools/xmcp/.venv/bin/python server.py
Restart=on-failure
RestartSec=5

# --- Hardening (L3) ---
# Baseline (also worked under --user):
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# Now enforceable as a SYSTEM unit (these silently failed in --user mode):
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectProc=invisible
ProcSubset=pid

[Install]
WantedBy=multi-user.target
